# templates/init.conf.erb
description "consular"

start on stopped rc RUNLEVEL=[2345]
respawn

exec /usr/bin/consular<%
  @final_opts.each do |k, v| %> <%= "--#{k} #{v}" %><% end %>
